# crontab
*/15 * * * *  python3 manage.py backup_all